# catcode.service
[Unit]
Description=CatCode - Feishu + Claude Code Bridge
After=network.target

[Service]
Type=simple
User=catcode
Group=catcode
WorkingDirectory=/opt/catcode
EnvironmentFile=/opt/catcode/.env
ExecStart=/root/.local/bin/uv run uvicorn catcode.main:app --host 0.0.0.0 --port ${PORT:-3000}
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
